# systemd/cryptogotchi-bluetooth-agent.service
[Unit]
Description=CryptoGotchi restricted Bluetooth pairing agent
After=bluetooth.service dbus.service
Requires=bluetooth.service

[Service]
Type=simple
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/local/lib/cryptogotchi/cryptogotchi-bluetooth-agent.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run

[Install]
WantedBy=multi-user.target
